# linux_service_configuration/incapsulaLogs.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
setuid $USER$
setgid $GROUP$
exec python $PYTHON_SCRIPT$
